# phab-sshd.sv.conf
[program:phab-sshd]
command=/usr/sbin/sshd -D -f /etc/phabricator-ssh/sshd_config.phabricator