# sshd.sv.conf
[program:sshd]
command=/usr/sbin/sshd -D -f /etc/ssh/sshd_config


/var/lib/docker/devicemapper/mnt/3d1ee141df449a1fe1db3d6a414d36d4d0262d502f439c542e1a7518a2cb51ee/rootfs
/var/lib/docker/devicemapper/mnt/3d1ee141df449a1fe1db3d6a414d36d4d0262d502f439c542e1a7518a2cb51ee/rootfs