# upgrade-phabricator.cron
30 3 * * 0 root /srv/phabricator/upgrade-phabricator.sh
